# root/etc/services.d/murmur/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /app/murmur/murmur.x86 \
	-ini /config/murmur.ini -fg
